# wellness_checker.service.template
[Unit]
Description=Run python script to ask about wellness

[Service]
Type=simple
ExecStart=/usr/bin/python3 PATH_TO_SCRIPT
