# wellness_checker.timer
[Unit]
Description=Run service that asks you about your wellness

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
